# ixr7250x3b/service/ixr7250x3b_platform_init.service
[Unit]
Description=Nokia-7250-IXR-X3b Platform Service
After=sysinit.target
Before=determine-reboot-cause.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ixr7250x3b_platform_init.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ixr7250x3b_platform_init

[Install]
WantedBy=multi-user.target
