# debian/sonic-platform-nokia-chassis.postinst
systemctl enable nokia-watchdog.service
systemctl start nokia-watchdog.service

systemctl enable openbdb.service
systemctl start openbdb.service

/usr/local/bin/nokia_install_ndk_service
