# debian/sonic-platform-nokia-ixr7220h6-64.postinst
#!/bin/sh
# postinst script for sonic-platform-nokia-IXR7220-H6-64
#
# see: dh_installdeb(1)

chmod a+x /usr/local/bin/h6_64_platform_init.sh
systemctl enable h6_64_platform_init.service
systemctl start h6_64_platform_init.service
chmod a+x /usr/local/bin/ports_notify.py
systemctl enable ports_notify.service
systemctl start --no-block ports_notify.service
